# systemd/auto-update.service
[Unit]
Description=AutoMaintain checker

[Service]
ExecStart=/etc/monitor-scripts/auto-update.sh
