# systemd/lora-pkt-fwd.timer
[Unit]
Description=Starts packet forwarder.

[Timer]
OnBootSec=20

[Install]
WantedBy=timers.target
