# systemd/pantherx-watchdog.service
[Unit]
Description=PantherX watchdog checker

[Service]
ExecStart=/etc/monitor-scripts/pantherx-watchdog.sh
