# systemd/server-detection.service
[Unit]
Description=Server detection checker

[Service]
ExecStart=/etc/monitor-scripts/server-detection.sh
