# systemd/timezone-config-check.service
[Unit]
Description=Timezone Config checker

[Service]
ExecStart=/etc/monitor-scripts/timezone-config.sh
