# systemd/update-dashboard-check.service
[Unit]
Description=Dashboard updater checker

[Service]
ExecStart=/etc/monitor-scripts/dashboard-update.sh
